# roles/desktop/files/70-wifi-powersave.rules
# Disable WiFi power management to prevent connection drops
# Applies to all wireless network interfaces
ACTION=="add", SUBSYSTEM=="net", KERNEL=="wlan*", RUN+="/usr/bin/iw dev %k set power_save off"